# supervisord.conf
[program:bot]
command = /home/dotcloud/env/bin/python /home/dotcloud/current/bin/bender
environment=PYTHONUNBUFFERED=true
