# supervisord.conf
[program:bot]
command = python /home/dotcloud/current/run.py
